# ansible/roles/kubernetes_node/files/etc/init/docker-bootstrap.conf
# A bootstrap Docker instance, which is used to start etcd & flannel

description     "Bootstrap Docker (used to start etcd and flanneld)"

start on starting docker
stop on runlevel [!2345]

exec dockerd -H unix:///var/run/docker-bootstrap.sock -p /var/run/docker-bootstrap.pid --iptables=false --ip-masq=false --bridge=none --graph=/var/lib/docker-bootstrap --exec-root=/var/run/docker-bootstrap
